# config/systemd/qcg-vif-monitor.path
[Unit]
Description=Monitor new VIF interfaces

[Path]
PathModified=/sys/class/net
MakeDirectory=yes

[Install]
WantedBy=multi-user.target
